# ubuntu.crontab
* * * * * /home/ubuntu/start_voltdb_if_needed.sh >> /voltdbdata/voltdbroot/log/cron.log
